# data/docker/supervisord.conf
[supervisord]
nodaemon=false

[program:litServer]
command=/home/lit-server.sh
autostart=true
autorestart=false
stdout_logfile=/logs/docker-server.out
stderr_logfile=/logs/docker-server.err
